# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-worker]
; Proses queue worker untuk notifikasi, deadline checks, dan job lainnya.
; --tries=3   : retry 3x sebelum job masuk failed_jobs
; --timeout=90: kill job yang berjalan > 90 detik (cegah zombie)
; --sleep=3   : polling interval saat queue kosong
; --max-jobs=500: restart worker setiap 500 job (cegah memory leak)
command=php /var/www/html/artisan queue:work --tries=3 --timeout=90 --sleep=3 --max-jobs=500
user=www-data
autostart=true
autorestart=true
priority=30
stopwaitsecs=120
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
